# systemd/update-all-remote-wn-clients.timer
[Unit]
Description=Periodic remote WN client update
After=syslog.target network-online.target
Wants=network-online.target

[Timer]
OnActiveSec=0s
OnUnitActiveSec=12h
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target

